# supervisor-master.conf
[supervisord]
nodaemon=true
user=root
loglevel=info

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err

[program:rocket-server]
command=/bin/bash -c "\
    cd /app && \
    echo '[master] Building server binary...' && \
    cargo build --release --bin server 2>&1 && \
    echo '[master] Build complete. Waiting for worker SSH...' && \
    until ssh -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null \
              -o ConnectTimeout=3 root@mpi-worker true 2>/dev/null; do \
        echo '[master] Worker SSH not ready yet, retrying in 2s...'; \
        sleep 2; \
    done && \
    echo '[master] Worker SSH ready. Starting MPI cluster...' && \
    mpirun --allow-run-as-root \
           --mca plm_rsh_args '-o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null' \
           --hostfile /etc/mpi-hostfile \
           -np 2 \
           /app/target/release/server"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/app.log
stderr_logfile=/var/log/app.err
stopasgroup=true
environment=OMPI_ALLOW_RUN_AS_ROOT=1,OMPI_ALLOW_RUN_AS_ROOT_CONFIRM=1,RUST_BACKTRACE=1
startretries=5
